# .config/autostart/discord.desktop
[Desktop Entry]
Exec=i3-msg 'workspace 4; exec discord'
Name=discord
Type=Application
