# .config/autostart/feh.desktop
[Desktop Entry]
Exec=/home/sizzleru/.local/bin/feh.sh
Name=feh
Type=Application
