# .config/autostart/ibus-daemon.desktop
[Desktop Entry]
Exec=bash -c 'ibus-daemon -drx'
Name=ibus-daemon
Type=Application
