# .config/autostart/polybar.desktop
[Desktop Entry]
Exec=/home/sizzleru/.config/polybar/launch.sh
Name=polybar
Type=Application
