# .config/autostart/streamdeck.desktop
[Desktop Entry]
Exec=bash -c 'streamdeck -n'
Name=streamdeck
Type=Application
