# .config/autostart/thunderbird.desktop
[Desktop Entry]
Exec=thunderbird
Name=thunderbird
Type=Application
