# .config/autostart/todoist.desktop
[Desktop Entry]
Exec=i3-msg 'workspace 3; exec todoist'
Name=todoist
Type=Application
