# .config/autostart/todoist.desktop
[Desktop Entry]
Exec=todoist
Name=todoist
Type=Application
